# mcdr/backup-remote-rsync/ten-bak.timer
[Unit]
Description=ten mc backup


[Timer]
OnBootSec=15m
OnUnitActiveSec=3h

[Install]
WantedBy=default.target
